# meta-siemens/recipes-example/example/files/button-led.service
[Unit]
Description=Button LED Control Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/button-led
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target